# build/etc/supervisord.conf
[supervisord]
logfile=/opt/esdb/kafka-agent/build/log/supervisord.log
pidfile=/opt/esdb/kafka-agent/build/log/supervisord.pid
childlogdir=/opt/esdb/kafka-agent/build/log/
nodaemon=true

[program:zookeeper]
command=/opt/esdb/kafka-agent/build/kafka/bin/zookeeper-server-start.sh /opt/esdb/kafka-agent/build/etc/zookeeper.properties
stdout_logfile=/opt/esdb/kafka-agent/build/log/zookeeper.log
stderr_logfile=NONE
redirect_stderr=true

[program:kafka]
command=/opt/esdb/kafka-agent/build/kafka/bin/kafka-server-start.sh /opt/esdb/kafka-agent/build/etc/server.properties
stdout_logfile=/opt/esdb/kafka-agent/build/log/zookeeper.log
stderr_logfile=NONE
redirect_stderr=true